# crontab
*/5 * * * * (cd /home/pi/source && /home/pi/source/mqtt_dh11.sh)