# crontab
* * * * * (cd /home/pi/source/temp_sensor && /home/pi/source/temp_sensor/mqtt_dh11.sh)
