# crontab
* * * * * (cd /home/pi/source/pi_temp_humidity_sensor && /home/pi/source/pi_temp_humidity_sensor/mqtt_dh11.sh)
